# platform/backup/survive-backup.service
[Unit]
Description=SURVIVE OS Backup Service (API)
After=network.target

[Service]
Type=simple
ExecStart=/opt/survive-os/venv/bin/uvicorn platform.backup.app.main:app --host 127.0.0.1 --port 8095
WorkingDirectory=/opt/survive-os
User=survive
Group=survive
Restart=on-failure
RestartSec=10
StandardOutput=journal
StandardError=journal
SyslogIdentifier=survive-backup

# Hardening
NoNewPrivileges=yes
ProtectSystem=strict
ReadWritePaths=/var/lib/survive /mnt/backup
ProtectHome=yes
PrivateTmp=yes

[Install]
WantedBy=multi-user.target
