# medical/concepts/survive-concepts.service
[Unit]
Description=SURVIVE OS Clinical Concept Dictionary
After=network.target
Wants=network.target

[Service]
Type=simple
User=survive
Group=survive
WorkingDirectory=/opt/survive/medical/concepts
ExecStart=/opt/survive/medical/venv/bin/uvicorn app.main:app --host 0.0.0.0 --port 8041
Restart=on-failure
RestartSec=5
StandardOutput=journal
StandardError=journal
SyslogIdentifier=survive-concepts

# Hardening
NoNewPrivileges=yes
ProtectSystem=strict
ReadWritePaths=/var/lib/survive/concepts
ProtectHome=yes
PrivateTmp=yes

[Install]
WantedBy=multi-user.target
