# platform/image-build/overlays/base/etc/systemd/system/survive-first-boot.service
[Unit]
Description=SURVIVE OS First Boot Initialization
ConditionPathExists=!/var/lib/survive/.first-boot-complete
After=network.target local-fs.target
Before=nginx.service

[Service]
Type=oneshot
ExecStart=/usr/lib/survive/first-boot.sh
RemainAfterExit=yes
StandardOutput=journal+console
StandardError=journal+console

[Install]
WantedBy=multi-user.target
